# webapp/scripts/stock.cron
#! /bin/bash

/var/www/stock/webapp/stock backup
