# scripts/mystocks-go.cron
#! /bin/bash

/var/www/mystocks-go/mystocks-go backup